# package/mtcdt/control/prerm
#!/bin/sh

/usr/sbin/update-rc.d -f helium_gateway remove


