# package/bobcat/control/prerm
#!/bin/sh

/usr/sbin/update-rc.d helium_gateway -f remove


