# package/resiot/control/postinst
#!/bin/sh

chmod +x /etc/helium_gateway/install_update

/usr/sbin/update-rc.d helium_gateway defaults 90
/etc/init.d/helium_gateway restart
